# agent/launchd/com.swil.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.swil.heartbeat.plist
  ========================
  macOS launchd 配置文件，用于保活并定期唤醒 heartbeat.sh。

  安装路径：~/Library/LaunchAgents/com.swil.heartbeat.plist
  （详细维护命令见 agent/MAINTENANCE.md）
-->
<plist version="1.0">
<dict>

  <!-- 唯一标识符，launchctl 命令使用这个 label -->
  <key>Label</key>
  <string>com.swil.heartbeat</string>

  <!-- 要执行的命令 -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/supwils/supwilsoft/swil/swil-social/agent/scripts/heartbeat.sh</string>
  </array>

  <!-- 工作目录 -->
  <key>WorkingDirectory</key>
  <string>/Users/supwils/supwilsoft/swil/swil-social/agent</string>

  <!-- 环境变量：确保 PATH 包含 claude CLI 和系统工具 -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Applications/cmux.app/Contents/Resources/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/supwils</string>
  </dict>

  <!-- 加载后立即启动一次 -->
  <key>RunAtLoad</key>
  <true/>

  <!-- 进程退出后自动重启（保活） -->
  <key>KeepAlive</key>
  <true/>

  <!-- 标准输出日志 -->
  <key>StandardOutPath</key>
  <string>/Users/supwils/supwilsoft/swil/swil-social/agent/logs/launchd-stdout.log</string>

  <!-- 标准错误日志 -->
  <key>StandardErrorPath</key>
  <string>/Users/supwils/supwilsoft/swil/swil-social/agent/logs/launchd-stderr.log</string>

  <!-- 节流保护：进程崩溃后至少等 30 秒再重启，避免崩溃循环 -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

</dict>
</plist>
